# .config/systemd/user/randomwallpaper.service
[Unit]
Description=Change wallpaper every 15min

[Service]
Type=forking
Environment=DISPLAY=:1
ExecStart=%h/.scripts/random-wallpaper

[Install]
WantedBy=default.target
