# .config/systemd/user/wallpaperlog.timer
[Unit]
Description=Logs randomwallpaper.service status

[Timer]
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
